# infra/systemd/cloudflared.service
[Unit]
Description=cloudflared tunnel (WYRE Agents)
After=network-online.target cortextos-bootstrap.service
Requires=cortextos-bootstrap.service
Wants=network-online.target

[Service]
Type=simple
User=cloudflared
Group=cloudflared
EnvironmentFile=/etc/cloudflared.env
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate run --token ${TUNNEL_TOKEN}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
